# bin/celery-worker-supervisor.conf
[program:celery-worker-default-supervisor]
environment=LSS_ENV=prod
command=/home/ubuntu/api_server/.env/bin/celery -A project.celery_config worker -l info -Q default_queue --concurrency=4 -O fair
directory=/home/ubuntu/api_server/
user=ubuntu
numprocs=1
stdout_logfile=/home/ubuntu/logs/celery-worker-notification.log
stderr_logfile=/home/ubuntu/logs/celery-worker-notification.log
autostart=true
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998
